# systemd/dgx-app.service
# Install:
#   sudo cp /home/echomind/dgx-local-app/systemd/dgx-app.service /etc/systemd/system/dgx-app.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now dgx-app
#
# Boot: uvicorn starts first, then (after a short delay) maybe_start_hotspot_stack.sh runs unless
# /tmp/dgx-hotspot-skip-autostart exists (after Disable/Safe stop until Enable or reboot).
# Then start_hotspot_stack or hotspot_enable if the configured marker is missing.
# echomind MUST have NOPASSWD + !requiretty for hotspot_enable.sh and start_hotspot_stack.sh
# (maybe_start invokes them via sudo) — run scripts/install_passwordless_sudo.sh
# First boot with no configured marker can run apt from hotspot_enable (may take minutes once).
# Leading "-" = hotspot failure does not fail the whole dgx-app unit.

[Unit]
Description=DGX App
After=network.target NetworkManager.service
Wants=network.target

[Service]
User=echomind
WorkingDirectory=/home/echomind/dgx-local-app/backend
Environment=DGX_APP_PORT=28734
ExecStart=/home/echomind/dgx-local-app/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 28734
# Cold boot: wlan driver / NM settle before hostapd (tune sleep if needed).
ExecStartPost=-/bin/sleep 5
ExecStartPost=-/bin/bash /home/echomind/dgx-local-app/scripts/maybe_start_hotspot_stack.sh
Restart=always

[Install]
WantedBy=multi-user.target
